# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

#[program:rails_app]
#command=/rails_start.sh